# ops/com.blob.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blob.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/red.g/CascadeProjects/Blob/.venv/bin/python</string>
        <string>-m</string>
        <string>blob</string>
        <string>loop</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/red.g/CascadeProjects/Blob</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/red.g/CascadeProjects/Blob/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/red.g/CascadeProjects/Blob/logs/launchd.err.log</string>
</dict>
</plist>
